# src/core/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock    ; (the path to the socket file)
chmod=0777                       ; socket file mode (default 0700)

[inet_http_server]
port = *:9001
username = %(ENV_SUPERVISOR_LOGIN)s
password = %(ENV_SUPERVISOR_PASSWORD)s
[supervisord]
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
user=root
serverurl=unix:///var/run/supervisor.sock


[program:main]
user=root
directory=/opt/catsearch/cat-backend
command=gunicorn -w %(ENV_WORKERS)s -k uvicorn.workers.UvicornWorker -b 0.0.0.0:8080 app.main:app --timeout %(ENV_WORKER_TIMEOUT)s
stdout_logfile=/var/log/catbackend_out.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10
startretries=300
stopwaitsecs=60
